# antivirus/clamscan.timer
[Unit]
Description=ClamAV virus scan

[Timer]
OnCalendar=*-*-1,15 0:00:00
Persistent=true

[Install]
WantedBy=timers.target
